# deploy/systemd/cortex-health.timer
[Unit]
Description=Probe CorTeX frontend health every 30s

[Timer]
OnBootSec=60
OnUnitActiveSec=30
AccuracySec=5s

[Install]
WantedBy=timers.target
