# deploy/systemd/cortex-health.service
[Unit]
Description=CorTeX frontend health probe (restarts a hung frontend)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cortex-healthcheck.sh
